# supervisor/puch-api.ini
[program:puch-api]
command=/srv/puch/venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8200 --workers 3 --timeout 120
directory=/srv/puch/backend
user=almalinux
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/puch-api.log
stderr_logfile=/var/log/supervisor/puch-api-error.log
environment=
    DJANGO_SETTINGS_MODULE="config.settings",
    DJANGO_DEBUG="False"
